# .config/macos/Library/LaunchAgents/com.user.spotlight.applications.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.spotlight.applications</string>

    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>-c</string>
        <string>~/.local/share/macos/spotlight-applications-import.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Applications</string>
        <string>~/Applications</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.user.spotlight.applications.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.user.spotlight.applications.err</string>
</dict>
</plist>
